# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/init-plugin-store/up
/etc/s6-overlay/scripts/setup-plugin-store.sh
