# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/svc-casa-mcp/finish
#!/command/with-contenv bashio
# 5.5 item 3: strip ANSI from bashio output for clean docker logs.
export BASHIO_LOG_NO_COLORS=true
export NO_COLOR=1

# svc-casa-mcp going down does NOT take the addon down — it's supportive.
# Engagement subprocesses lose MCP access until the supervisor respawns
# the service (next 5s under s6-rc default policy). casa-main keeps
# serving its public-port fallback for any pre-v0.14.0 workspaces.
if [ "${1}" -ne 0 ] && [ "${1}" -ne 256 ]; then
    bashio::log.warning "svc_casa_mcp exited with code ${1} — supervisor will respawn"
fi
