# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/svc-casa/finish
#!/command/with-contenv bashio

if [ "${1}" -ne 0 ] && [ "${1}" -ne 256 ]; then
    bashio::log.error "Casa crashed with exit code ${1}"
    bashio::addon.stop
fi
