# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/svc-casa/run
#!/command/with-contenv bashio
# 5.5 item 3: strip ANSI from bashio output for clean docker logs.
export BASHIO_LOG_NO_COLORS=true
export NO_COLOR=1

# bashio::config returns "null" for empty optional strings; normalize to empty
_public_url="$(bashio::config 'public_url')"
if [ "$_public_url" = "null" ] || [ -z "$_public_url" ]; then
    export PUBLIC_URL=""
else
    export PUBLIC_URL="$_public_url"
fi
export CLAUDE_CODE_OAUTH_TOKEN="$(bashio::config 'claude_oauth_token')"
export OP_SERVICE_ACCOUNT_TOKEN="$(bashio::config 'onepassword_service_account_token')"
export ONEPASSWORD_DEFAULT_VAULT="$(bashio::config 'onepassword_default_vault')"
export HONCHO_API_KEY="$(bashio::config 'honcho_api_key')"
export HONCHO_API_URL="$(bashio::config 'honcho_api_url')"
export TELEGRAM_BOT_TOKEN="$(bashio::config 'telegram_bot_token')"
export TELEGRAM_CHAT_ID="$(bashio::config 'telegram_chat_id')"
_tg_transport="$(bashio::config 'telegram_transport')"
export TELEGRAM_TRANSPORT="${_tg_transport:-polling}"
export TELEGRAM_DELIVERY_MODE="$(bashio::config 'telegram_delivery_mode')"
export WEBHOOK_SECRET="$(bashio::config 'webhook_secret')"
export PRIMARY_AGENT_NAME="$(bashio::config 'primary_agent_name')"
export VOICE_AGENT_NAME="$(bashio::config 'voice_agent_name')"
export PRIMARY_AGENT_MODEL="$(bashio::config 'primary_agent_model')"
export VOICE_AGENT_MODEL="$(bashio::config 'voice_agent_model')"
export SUBAGENT_MODEL="$(bashio::config 'subagent_model')"
export ENABLE_TERMINAL="$(bashio::config 'enable_terminal')"
export CASA_VERSION="$(bashio::addon.version)"

# 3.2: fastembed model cache — must match the path used by
# setup-configs.sh pre-warm so the runtime finds the downloaded model.
export FASTEMBED_CACHE_PATH="/data/fastembed"

# 3.2.1: scope classifier threshold — promoted from env-only default to
# first-class addon option so Builder can tune it via supervisor.addon_options_set.
# Runtime read at casa_core.py:427 is unchanged; only the source of the
# env var changes (was: fallback 0.35; now: bashio::config value).
# "null"-normalize matches the PUBLIC_URL pattern: if a user clears the
# key via YAML UI, leave the var unset so casa_core.py:427 falls back to "0.35"
# rather than attempting float("null").
_scope_threshold="$(bashio::config 'scope_threshold')"
if [ "$_scope_threshold" != "null" ] && [ -n "$_scope_threshold" ]; then
    export CASA_SCOPE_THRESHOLD="$_scope_threshold"
fi

# 3.3: scheduler timezone — optional; timekeeping.resolve_tz() falls back
# to the TZ env var set by HA OS, then Europe/Amsterdam.
_casa_tz="$(bashio::config 'casa_tz')"
if [ "$_casa_tz" != "null" ] && [ -n "$_casa_tz" ]; then
    export CASA_TZ="$_casa_tz"
fi

# SUPERVISOR_TOKEN is injected automatically by the HA Supervisor.

exec python3 /opt/casa/casa_core.py
