# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/svc-casa/run
#!/command/with-contenv bashio
# 5.5 item 3: strip ANSI from bashio output for clean docker logs.
export BASHIO_LOG_NO_COLORS=true
export NO_COLOR=1

# bashio::config returns "null" for empty optional strings; normalize to empty
_public_url="$(bashio::config 'public_url')"
if [ "$_public_url" = "null" ] || [ -z "$_public_url" ]; then
    export PUBLIC_URL=""
else
    export PUBLIC_URL="$_public_url"
fi
export CLAUDE_CODE_OAUTH_TOKEN="$(bashio::config 'claude_oauth_token')"
export OP_SERVICE_ACCOUNT_TOKEN="$(bashio::config 'onepassword_service_account_token')"
export ONEPASSWORD_DEFAULT_VAULT="$(bashio::config 'onepassword_default_vault')"
_hindsight_url="$(bashio::config 'hindsight_api_url')"
if [ "$_hindsight_url" != "null" ] && [ -n "$_hindsight_url" ]; then
    export HINDSIGHT_API_URL="$_hindsight_url"
    # Setting hindsight_api_url is the single toggle for long-term memory:
    # casa_core requires MEMORY_BACKEND=hindsight (anything else → noop), and
    # nothing else sets it, so derive it here. Respect an explicit override.
    export MEMORY_BACKEND="${MEMORY_BACKEND:-hindsight}"
fi
export TELEGRAM_BOT_TOKEN="$(bashio::config 'telegram_bot_token')"
export TELEGRAM_CHAT_ID="$(bashio::config 'telegram_chat_id')"
_tg_transport="$(bashio::config 'telegram_transport')"
export TELEGRAM_TRANSPORT="${_tg_transport:-polling}"
export TELEGRAM_DELIVERY_MODE="$(bashio::config 'telegram_delivery_mode')"
export TELEGRAM_ENGAGEMENT_SUPERGROUP_ID="$(bashio::config 'telegram_engagement_supergroup_id')"
export WEBHOOK_SECRET="$(bashio::config 'webhook_secret')"
export PRIMARY_AGENT_NAME="$(bashio::config 'primary_agent_name')"
export VOICE_AGENT_NAME="$(bashio::config 'voice_agent_name')"
export PRIMARY_AGENT_MODEL="$(bashio::config 'primary_agent_model')"
export VOICE_AGENT_MODEL="$(bashio::config 'voice_agent_model')"
export SUBAGENT_MODEL="$(bashio::config 'subagent_model')"
export ENABLE_TERMINAL="$(bashio::config 'enable_terminal')"
export CASA_VERSION="$(bashio::addon.version)"

# 3.3: scheduler timezone — optional; timekeeping.resolve_tz() falls back
# to the TZ env var set by HA OS, then Europe/Amsterdam.
_casa_tz="$(bashio::config 'casa_tz')"
if [ "$_casa_tz" != "null" ] && [ -n "$_casa_tz" ]; then
    export CASA_TZ="$_casa_tz"
fi

# v0.18.1: operator-facing log level. install_logging() at
# casa_core.py:743 reads LOG_LEVEL env when set; falls back to INFO
# when unset/null/empty.
_log_level="$(bashio::config 'log_level')"
if [ "$_log_level" != "null" ] && [ -n "$_log_level" ]; then
    export LOG_LEVEL="$_log_level"
fi

# SUPERVISOR_TOKEN is injected automatically by the HA Supervisor.

exec python3 /opt/casa/casa_core.py
