# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/svc-casa/run
#!/command/with-contenv bashio

export CLAUDE_CODE_OAUTH_TOKEN="$(bashio::config 'claude_oauth_token')"
export HONCHO_API_KEY="$(bashio::config 'honcho_api_key')"
export HONCHO_API_URL="$(bashio::config 'honcho_api_url')"
export TELEGRAM_BOT_TOKEN="$(bashio::config 'telegram_bot_token')"
export TELEGRAM_CHAT_ID="$(bashio::config 'telegram_chat_id')"
export TELEGRAM_WEBHOOK_URL="$(bashio::config 'telegram_webhook_url')"
export TELEGRAM_DELIVERY_MODE="$(bashio::config 'telegram_delivery_mode')"
export WEBHOOK_SECRET="$(bashio::config 'webhook_secret')"
export PRIMARY_AGENT_NAME="$(bashio::config 'primary_agent_name')"
export VOICE_AGENT_NAME="$(bashio::config 'voice_agent_name')"
export PRIMARY_AGENT_MODEL="$(bashio::config 'primary_agent_model')"
export VOICE_AGENT_MODEL="$(bashio::config 'voice_agent_model')"
export SUBAGENT_MODEL="$(bashio::config 'subagent_model')"
export ENABLE_TERMINAL="$(bashio::config 'enable_terminal')"
export CASA_VERSION="$(bashio::addon.version)"

# SUPERVISOR_TOKEN is injected automatically by the HA Supervisor.

exec python3 /opt/casa/casa_core.py
