# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/svc-casa/run
#!/command/with-contenv bashio
# 5.5 item 3: strip ANSI from bashio output for clean docker logs.
export BASHIO_LOG_NO_COLORS=true
export NO_COLOR=1

# bashio::config returns "null" for empty optional strings; normalize to empty
_public_url="$(bashio::config 'public_url')"
if [ "$_public_url" = "null" ] || [ -z "$_public_url" ]; then
    export PUBLIC_URL=""
else
    export PUBLIC_URL="$_public_url"
fi
export CLAUDE_CODE_OAUTH_TOKEN="$(bashio::config 'claude_oauth_token')"
export OP_SERVICE_ACCOUNT_TOKEN="$(bashio::config 'onepassword_service_account_token')"
export ONEPASSWORD_DEFAULT_VAULT="$(bashio::config 'onepassword_default_vault')"
_hindsight_url="$(bashio::config 'hindsight_api_url')"
if [ "$_hindsight_url" != "null" ] && [ -n "$_hindsight_url" ]; then
    export HINDSIGHT_API_URL="$_hindsight_url"
    # Setting hindsight_api_url is the single toggle for long-term memory:
    # casa_core requires MEMORY_BACKEND=hindsight (anything else → noop), and
    # nothing else sets it, so derive it here. Respect an explicit override.
    export MEMORY_BACKEND="${MEMORY_BACKEND:-hindsight}"
fi
export TELEGRAM_BOT_TOKEN="$(bashio::config 'telegram_bot_token')"
export TELEGRAM_CHAT_ID="$(bashio::config 'telegram_chat_id')"
_tg_transport="$(bashio::config 'telegram_transport')"
export TELEGRAM_TRANSPORT="${_tg_transport:-polling}"
export TELEGRAM_DELIVERY_MODE="$(bashio::config 'telegram_delivery_mode')"
export TELEGRAM_RICH_TEXT="$(bashio::config 'telegram_rich_text')"
export TELEGRAM_ENGAGEMENT_SUPERGROUP_ID="$(bashio::config 'telegram_engagement_supergroup_id')"
# telegram_bot_api_base is optional (schema str?); bashio returns "null" when unset.
_tg_api_base="$(bashio::config 'telegram_bot_api_base')"
if [ "$_tg_api_base" != "null" ] && [ -n "$_tg_api_base" ]; then
    export TELEGRAM_BOT_API_BASE="$_tg_api_base"
fi
# webhook_auth_enabled gates whether the webhook secret check is active;
# only export WEBHOOK_SECRET when the toggle is on, so setting it to false
# actually disables auth even if webhook_secret still holds a value.
if bashio::config.true 'webhook_auth_enabled'; then
    export WEBHOOK_SECRET="$(bashio::config 'webhook_secret')"
fi
export PRIMARY_AGENT_MODEL="$(bashio::config 'primary_agent_model')"
export VOICE_AGENT_MODEL="$(bashio::config 'voice_agent_model')"
export ENABLE_TERMINAL="$(bashio::config 'enable_terminal')"
# sdk_client_pool: operator-reachable kill-switch for the warm SDK client
# pool (rollback path). pool_enabled() in sdk_client_pool.py treats
# off/0/false as disabled and anything else as enabled, so exporting
# bashio's raw "true"/"false" string works the same way as ENABLE_TERMINAL.
export SDK_CLIENT_POOL="$(bashio::config 'sdk_client_pool')"
# Tina's eager Home Assistant facade is default-on with an operator-reachable
# raw-MCP fallback for diagnosis.
export TINA_HA_FACADE_ENABLED="$(bashio::config 'tina_ha_facade_enabled')"
export CASA_VERSION="$(bashio::addon.version)"

# 3.3: scheduler timezone — optional; timekeeping.resolve_tz() falls back
# to the TZ env var set by HA OS, then Europe/Amsterdam.
_casa_tz="$(bashio::config 'casa_tz')"
if [ "$_casa_tz" != "null" ] && [ -n "$_casa_tz" ]; then
    export CASA_TZ="$_casa_tz"
fi

# v0.18.1: operator-facing log level. install_logging() at
# casa_core.py:743 reads LOG_LEVEL env when set; falls back to INFO
# when unset/null/empty.
_log_level="$(bashio::config 'log_level')"
if [ "$_log_level" != "null" ] && [ -n "$_log_level" ]; then
    export LOG_LEVEL="$_log_level"
fi

# D-4 (v0.69.0): stale-engagement reap TTL in days; 0 disables the reap.
# tools._engagement_reap_days() falls back to 7 when unset/invalid.
_reap_days="$(bashio::config 'engagement_reap_days')"
if [ "$_reap_days" != "null" ] && [ -n "$_reap_days" ]; then
    export ENGAGEMENT_REAP_DAYS="$_reap_days"
fi

# Task 7 (spec A4/§4.6): voice turn budget + specialist concurrency/cost-alert
# options. All three ship with a default in the options block above, so
# bashio never returns "null" here (unlike the optional str?/int? options
# above) — export unconditionally, matching TELEGRAM_ENGAGEMENT_SUPERGROUP_ID.
# _voice_turn_budget_s() in channels/voice/channel.py and the
# specialist_max_concurrency/specialist_cost_alert_threshold reads in
# casa_core.py still clamp/validate on top of this (e.g. the 27s hard cap) —
# this export only feeds the env var, it does not itself enforce the cap.
export VOICE_TURN_BUDGET_SECONDS="$(bashio::config 'voice_turn_budget_seconds')"
export VOICE_ROUTE_FRESHNESS_SECONDS="$(bashio::config 'voice_route_freshness_seconds')"
export VOICE_JOB_DELIVERY_TTL_SECONDS="$(bashio::config 'voice_job_delivery_ttl_seconds')"
export VOICE_JOB_ROUTE_CAP="$(bashio::config 'voice_job_route_cap')"
export SPECIALIST_MAX_CONCURRENCY="$(bashio::config 'specialist_max_concurrency')"
export SPECIALIST_COST_ALERT_THRESHOLD="$(bashio::config 'specialist_cost_alert_threshold')"

# SUPERVISOR_TOKEN is injected automatically by the HA Supervisor.

exec /opt/casa/venv/bin/python3 /opt/casa/casa_core.py
