# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/finish
#!/command/with-contenv bashio

bashio::log.warning "Nginx exited with code ${1}"
