# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/init-setup-configs/up
/etc/s6-overlay/scripts/setup-configs.sh
