# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/svc-ttyd/finish
#!/command/with-contenv bashio
# 5.5 item 3: strip ANSI from bashio output for clean docker logs.
export BASHIO_LOG_NO_COLORS=true
export NO_COLOR=1

if ! bashio::config.true 'enable_terminal'; then
    exit 256
fi

bashio::log.warning "ttyd exited with code ${1}"
