# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/svc-ttyd/finish
#!/command/with-contenv bashio
# 5.5 item 3: strip ANSI from bashio output for clean docker logs.
export BASHIO_LOG_NO_COLORS=true
export NO_COLOR=1

if ! bashio::config.true 'enable_terminal'; then
    exit 256
fi

# Code 0 = clean stop (s6 told it to). Code 256 = s6 "do-not-restart" sentinel.
# Anything else means ttyd crashed and the supervisor will respawn it.
if [ "${1}" -ne 0 ] && [ "${1}" -ne 256 ]; then
    bashio::log.warning "ttyd exited with code ${1} — supervisor will respawn"
fi
