# casa-agent/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
# /etc/services.d/nginx/finish

bashio::log.warning "Nginx exited with code ${1}"
