# casa-agent/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# /etc/services.d/nginx/run

# Wait for aiohttp to be ready before proxying (prevents 502 on startup)
bashio::net.wait_for 8099 127.0.0.1 30

exec nginx -g 'daemon off;'
