# casa-agent/rootfs/etc/s6-overlay/s6-rc.d/init-validate-config/up
/etc/s6-overlay/scripts/validate-config.sh
